# deploy/supervisor/tg-marketing-backend.conf
[program:tg-marketing-backend]
directory=/var/www/tg-marketing/backend
command=/var/www/tg-marketing/backend/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --proxy-headers
user=www-data
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/tg-marketing/backend.log
stderr_logfile=/var/log/tg-marketing/backend-error.log
environment=PATH="/var/www/tg-marketing/backend/.venv/bin"
